# files/supervisor.conf
{%- set service = salt['pillar.get']('supervisor:server:service:'+service_name) %}

[program:{{ service_name }}]
user=octoprint
directory=/srv/octoprint
environment=PATH="/srv/octoprint/env/bin"
command=/srv/octoprint/env/bin/server/start
stdout_logfile=/var/log/octoprint/{{ service.name }}_app.log
stderr_logfile=/var/log/octoprint/{{ service.name }}_err.log
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
